# rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv sh
set -e

crond -f